# extras/supervisor/aspyrobotmx.ini
[program:robot-server]
user=blctl
directory=/srv/aspyrobotmx
command=/srv/aspyrobotmx/.venv/bin/pyrobotmxserver --config /srv/aspyrobotmx/config.json SR08ID01ROB01
environment=PYTHONUNBUFFERED=1,BEAMLINE=Cage1,PYEPICS_LIBCA=/epics/base/lib/linux-x86_64/libca.so
stopasgroup=true
